# salt/salt/gunicorn/gunicorn_reloader.path
[Path]
PathChanged=/opt/tildes/static/css/site-icons.css

[Install]
WantedBy=multi-user.target
